# config/stirling-pdf/quadlet/default.container
[Unit]
Description=Stirling-PDF
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/stirlingtools/stirling-pdf:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=stirling-pdf
PublishPort=127.0.0.1:${host_port}:8080
#Network=tgdb.network

EnvironmentFile=${instance_dir}/.env

# 注意：若你不需要自訂 OCR 語言/資料，建議不要掛載 /usr/share/tessdata，避免覆蓋容器內建資料。
Volume=${instance_dir}/configs:/configs
Volume=${instance_dir}/customFiles:/customFiles
Volume=${instance_dir}/logs:/logs
Volume=${instance_dir}/pipeline:/pipeline

User=0:0

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
