# config/home-assistant/quadlet/default.container
[Unit]
Description=home-assistant
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/home-assistant/home-assistant:stable
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=home-assistant

EnvironmentFile=${instance_dir}/.env

# Home Assistant 常見部署會需要 host network 才能做區網/設備探索（mDNS/SSDP/UDP broadcast 等）
Network=host

Volume=${instance_dir}/config:/config

# --- 上游常見建議（風險較高；本範本預設已啟用，請只在可信任環境使用）---
PodmanArgs=--cgroupns=host
PodmanArgs=--privileged
Volume=/dev:/dev
Volume=/run/udev:/run/udev:ro
Volume=/etc/localtime:/etc/localtime:ro
Volume=/run/dbus:/run/dbus:ro

User=0:0

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false

