# config/twenty/quadlet/default.container
[Unit]
Description=twenty server
After=network-online.target
Wants=network-online.target

After=${container_name}-db.service container-${container_name}-db.service ${container_name}-redis.service container-${container_name}-redis.service
Wants=${container_name}-db.service container-${container_name}-db.service ${container_name}-redis.service container-${container_name}-redis.service

[Container]
Image=docker.io/twentycrm/twenty:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=twenty
UserNS=keep-id

Pod=${container_name}.pod

EnvironmentFile=${instance_dir}/.env

Volume=${instance_dir}/local-storage:/app/packages/twenty-server/.local-storage

[Service]
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 90 ]; do /usr/bin/podman exec ${container_name}-db pg_isready -h 127.0.0.1 -p 5432 -U ${db_user} >/dev/null 2>&1 && /usr/bin/podman exec ${container_name}-redis redis-cli ping >/dev/null 2>&1 && exit 0; i=$((i+1)); sleep 1; done; echo "PostgreSQL/Redis е°љжњЄе°±з·’пјЊеЏ–ж¶€е•џе‹• Twenty" >&2; exit 1'
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
