# config/alloy/quadlet/default.container
[Unit]
Description=Grafana Alloy
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/grafana/alloy:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=alloy

Environment=TZ=Asia/Taipei

PublishPort=127.0.0.1:${host_port}:12345

# 取用 host PID namespace，讓 /proc 類指標更貼近宿主機（對齊 node-exporter 容器做法）
PodmanArgs=--pid=host

Exec=run --server.http.listen-addr=0.0.0.0:12345 --storage.path=/var/lib/alloy/data /etc/alloy/alloy.config

Volume=${instance_dir}/alloy.config:/etc/alloy/alloy.config
Volume=${instance_dir}/data:/var/lib/alloy/data

User=0:0

Volume=/run/udev:/run/udev:ro
Volume=/proc:/host/proc:ro,rslave
Volume=/sys:/host/sys:ro,rslave
Volume=/:/host:ro,rslave

Volume=/var/log:/var/log:ro
Volume=/etc/machine-id:/etc/machine-id:ro
Volume=/run/log/journal:/run/log/journal:ro
Volume=/var/log/journal:/var/log/journal:ro

NoNewPrivileges=true

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
