# config/scrutiny/quadlet/default.container
[Unit]
Description=Scrutiny
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/analogj/scrutiny:v0.8.6-omnibus
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=scrutiny
#Network=tgdb.network

PublishPort=127.0.0.1:${host_port}:8080/tcp
#PublishPort=127.0.0.1:18086:8086/tcp

EnvironmentFile=${instance_dir}/.env

Volume=${instance_dir}/config:/opt/scrutiny/config
Volume=${instance_dir}/influxdb:/opt/scrutiny/influxdb
Volume=/run/udev:/run/udev:ro
Volume=/dev:/dev

# 上游文件常見需求：/run/udev、SYS_RAWIO；若有 NVMe 還可能需要 SYS_ADMIN。
# TGDB 這裡改採較保守但權限更高的方式，先盡量讓 collector 能正常讀到各種裝置。
PodmanArgs=--privileged

User=0:0

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
