# config/crucix/quadlet/default.container
[Unit]
Description=Crucix
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/calesthio/crucix:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=crucix

EnvironmentFile=${instance_dir}/.env

Volume=${instance_dir}/runs:/app/runs

PublishPort=127.0.0.1:${host_port}:3117/tcp

User=0
UserNS=keep-id

NoNewPrivileges=true

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
