# config/kestra/quadlet/default.container
[Unit]
Description=Kestra
After=network-online.target
Wants=network-online.target

After=${container_name}-postgres.service container-${container_name}-postgres.service
Wants=${container_name}-postgres.service container-${container_name}-postgres.service

[Container]
Image=docker.io/kestra/kestra:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=kestra

Pod=${container_name}.pod

Exec=server standalone --config /etc/kestra/application.yml

Volume=${instance_dir}/storage:/app/storage
Volume=${instance_dir}/tmp_wd:/tmp/kestra-wd
Volume=${instance_dir}/application.yml:/etc/kestra/application.yml:ro
Volume=${podman_sock_host_path}:/var/run/docker.sock

UserNS=keep-id
User=0

[Service]
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 90 ]; do /usr/bin/podman exec ${container_name}-postgres pg_isready -h 127.0.0.1 -p 5432 -U ${db_user} -d kestra >/dev/null 2>&1 && exit 0; i=$((i+1)); sleep 1; done; echo "PostgreSQL е°љжњЄе°±з·’пјЊеЏ–ж¶€е•џе‹• Kestra" >&2; exit 1'
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
