# config/kopia/quadlet/default.container
[Unit]
Description=Kopia（快照備份）
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/kopia/kopia:20260223.0.231822
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=kopia

#Network=tgdb.network

PublishPort=127.0.0.1:${host_port}:51515/tcp

EnvironmentFile=${instance_dir}/.env

Volume=${instance_dir}/config:/app/config
Volume=${instance_dir}/cache:/app/cache
Volume=${instance_dir}/logs:/app/logs
Volume=${volume_dir}:/repository
Volume=${backup_root}:/data:ro

Entrypoint=["/bin/sh","-c"]

Exec="set -eu; mkdir -p /app/rclone; if [ -f \"/data/${tgdb_name}/rclone.conf\" ]; then ln -sf \"/data/${tgdb_name}/rclone.conf\" /app/rclone/rclone.conf; fi; exec kopia server start --address=http://0.0.0.0:51515 --server-username=\"$KOPIA_SERVER_USERNAME\" --server-password=\"$KOPIA_SERVER_PASSWORD\" --insecure"

[Service]
Restart=on-failure
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
