# config/anything-llm/quadlet/default.container
[Unit]
Description=AnythingLLM
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/mintplexlabs/anythingllm:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=anything-llm

EnvironmentFile=${instance_dir}/.env

Volume=${instance_dir}/.env:/app/server/.env:ro
Volume=${instance_dir}/storage:/app/server/storage

PublishPort=127.0.0.1:${host_port}:3001/tcp
UserNS=keep-id:uid=0,gid=0

AddCapability=SYS_ADMIN

HealthCmd=sh -c 'curl -fsS http://127.0.0.1:3001/api/ping >/dev/null || exit 1'
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
