# config/kan/quadlet/default2.container
[Unit]
Description=Kan migration
After=network-online.target
Wants=network-online.target

After=${container_name}-postgres.service container-${container_name}-postgres.service
Wants=${container_name}-postgres.service container-${container_name}-postgres.service

[Container]
Image=ghcr.io/kanbn/kan-migrate:latest
#AutoUpdate=registry
ContainerName=${container_name}-migrate
Label=app=kan

Pod=${container_name}.pod

EnvironmentFile=${instance_dir}/.env

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 120 ]; do /usr/bin/podman exec ${container_name}-postgres pg_isready -h 127.0.0.1 -p 5432 -U ${db_user} -d kan >/dev/null 2>&1 && exit 0; i=$((i+1)); sleep 1; done; echo "PostgreSQL 尚未就緒，取消執行 Kan migration" >&2; exit 1'
Restart=on-failure
RestartSec=5s
StartLimitBurst=30

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
