# config/ackee/quadlet/default.container
[Unit]
Description=Ackee
After=network-online.target
Wants=network-online.target

After=${container_name}-mongo.service container-${container_name}-mongo.service
Wants=${container_name}-mongo.service container-${container_name}-mongo.service

[Container]
Image=docker.io/electerious/ackee:latest
#AutoUpdate=registry
ContainerName=${container_name}
Label=app=ackee

Pod=${container_name}.pod

EnvironmentFile=${instance_dir}/.env

[Service]
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 90 ]; do /usr/bin/podman exec ${container_name}-mongo mongosh --quiet --eval "db.adminCommand({ ping: 1 })" >/dev/null 2>&1 && exit 0; i=$((i+1)); sleep 1; done; echo "MongoDB е°љжњЄе°±з·’пјЊеЏ–ж¶€е•џе‹• Ackee" >&2; exit 1'
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
