# config/paperclip/quadlet/default.container
[Unit]
Description=Paperclip
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/paperclipai/paperclip:latest
ContainerName=${container_name}
Label=app=paperclip

EnvironmentFile=${instance_dir}/.env
PublishPort=127.0.0.1:${host_port}:3100/tcp

Volume=${instance_dir}/data:/paperclip

UserNS=keep-id
User=0

NoNewPrivileges=true

[Service]
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
