# config/misskey/quadlet/default.container
[Unit]
Description=Misskey
After=network-online.target
Wants=network-online.target

After=${container_name}-postgres.service container-${container_name}-postgres.service ${container_name}-redis.service container-${container_name}-redis.service
Wants=${container_name}-postgres.service container-${container_name}-postgres.service ${container_name}-redis.service container-${container_name}-redis.service

[Container]
Image=localhost/tgdb/misskey:${container_name}
ContainerName=${container_name}
Label=app=misskey

Pod=${container_name}.pod

UserNS=keep-id
User=0

Environment=NODE_ENV=production

Volume=${instance_dir}/.config:/misskey/.config:ro
Volume=${instance_dir}/files:/misskey/files

NoNewPrivileges=true

[Service]
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 120 ]; do /usr/bin/podman exec ${container_name}-postgres pg_isready -h 127.0.0.1 -p 5432 -U ${db_user} >/dev/null 2>&1 && /usr/bin/podman exec ${container_name}-redis redis-cli ping >/dev/null 2>&1 && exit 0; i=$((i+1)); sleep 1; done; echo "PostgreSQL/Redis е°љжњЄе°±з·’пјЊеЏ–ж¶€е•џе‹• Misskey" >&2; exit 1'
Restart=always
RestartSec=15s
StartLimitBurst=3

[Install]
WantedBy=default.target

[Quadlet]
DefaultDependencies=false
